# orangepi/systemd/orangepi-vision.service
[Unit]
Description=Orange Pi Vision Service
After=network.target photonvision.service
Wants=network-online.target
Requires=photonvision.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/orangepi-vision
ExecStart=/opt/orangepi-vision/venv/bin/python /opt/orangepi-vision/vision_service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Resource limits
Nice=-10
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50

[Install]
WantedBy=multi-user.target
